# deploy/dev/kraken-backend-dev.service
[Unit]
Description=Kraken Backend (Dev)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/kraken-chatbot-dev/backend
EnvironmentFile=/home/ubuntu/kraken-chatbot-dev/backend/.env
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn src.kestrel_backend.main:app --host 127.0.0.1 --port 8006
Restart=always
RestartSec=5
SyslogIdentifier=kraken-backend-dev

[Install]
WantedBy=multi-user.target
